# docker/root/etc/services.d/altstore_proxy/run
#!/usr/bin/with-contenv bash

exec altstore_proxy --port=8080 --cache=/cache --baseurl=$BASEURL --repos=$REPOS
